# linuxzones.service
[Unit]
Description=LinuxZones window zone manager
Documentation=https://github.com/v6belung/LinuxZones
After=graphical-session.target
PartOf=graphical-session.target

[Service]
Type=simple
ExecStart=%h/.local/bin/linuxzones
Restart=on-failure
RestartSec=5
# Most display managers propagate DISPLAY/XAUTHORITY to the systemd user daemon
# via dbus-update-activation-environment.  If linuxzones fails to start at login,
# run: systemctl --user show-environment | grep DISPLAY
# and see README.md → Troubleshooting for the manual fix.
PassEnvironment=DISPLAY XAUTHORITY WAYLAND_DISPLAY

[Install]
WantedBy=graphical-session.target
